# crontab
*/2 * * * * sh /sensu-client-removal.sh 
